# jetson_setup/jetson-deployment.service
[Unit]
Description=fprime-scales JetsonDeployment Flight Software
# Wait for network (needed to connect to the IMX hub)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace 'jetson' with the actual username on the Jetson
User=<Your User Name Here>
WorkingDirectory=~/fprime-scales-ref

ExecStart=~/fprime-scales-ref/jetson-startup.sh

# Restart on crash, but not on clean exit (exit 0)
Restart=on-failure
RestartSec=5

# Give the network and fprime-gds time to be ready before retrying hard failures
StartLimitIntervalSec=120
StartLimitBurst=5

# Log stdout/stderr to the journal (view with: journalctl -u jetson-deployment -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

